# systemd-user/jarvis-callprep.service
[Unit]
Description=JARVIS Call Prep
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/clawd/clawd/projects/call-prep
ExecStart=/usr/bin/python3 /home/clawd/clawd/projects/call-prep/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
