# launchagents/com.jarvis.golden-hour-scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jarvis.golden-hour-scout</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/python@3.14/bin/python3.14</string>
    <string>app.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/jarvis/clawd/projects/golden-hour-scout</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/golden-hour-scout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/golden-hour-scout-error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
